# 20240430/1.img/usr/lib/systemd/system/makerbase-time-monitor.service
[Unit]
Description=Timezone Monitor Service
After=network.target

[Service]
ExecStart=/root/set-time.sh

[Install]
WantedBy=multi-user.target
