# 20240430/1.img/usr/lib/systemd/system/makerbase-webcam.service
[Unit]
Description=AutoSetWebcamId
After=crowsnest.service


[Service]
Type=oneshot
ExecStart=bash /root/auto-webcam-id.sh

[Install]
WantedBy=multi-user.target



