# 20240430/1.img/var/lib/dpkg/info/giga-client.list
/.
/home
/home/mks
/home/mks/Desktop
/home/mks/Desktop/myfile
/home/mks/Desktop/myfile/znp
/home/mks/Desktop/myfile/znp/znp_tjc_klipper
/home/mks/Desktop/myfile/znp/znp_tjc_klipper/build
/home/mks/Desktop/myfile/znp/znp_tjc_klipper/build/io
/home/mks/Desktop/myfile/znp/znp_tjc_klipper/build/src
/home/mks/Desktop/myfile/znp/znp_tjc_klipper/build/src/libsrc.a
/home/mks/Desktop/myfile/znp/znp_tjc_klipper/build/start.sh
/home/mks/Desktop/myfile/znp/znp_tjc_klipper/build/znp_tjc_klipper
/home/mks/Desktop/myfile/znp/znp_tjc_klipper/c_helper.so
/home/mks/Desktop/myfile/znp/znp_tjc_klipper/canbus_query.py
/home/mks/Desktop/myfile/znp/znp_tjc_klipper/elegoo_plr.ini
/home/mks/Desktop/myfile/znp/znp_tjc_klipper/elegoo_plr.ini.bak
/home/mks/Desktop/myfile/znp/znp_tjc_klipper/elegoo_plr_data.ini
/home/mks/Desktop/myfile/znp/znp_tjc_klipper/elegoo_plr_data.ini.bak
/home/mks/Desktop/myfile/znp/znp_tjc_klipper/elegoo_plr_dataA.ini
/home/mks/Desktop/myfile/znp/znp_tjc_klipper/elegoo_plr_dataB.ini
/home/mks/Desktop/myfile/znp/znp_tjc_klipper/reset_canbus_id.py
/home/mks/Desktop/myfile/znp/znp_tjc_klipper/sshd_config_unlock
/home/mks/clear_plr.sh
/home/mks/klipper
/home/mks/klipper/klippy
/home/mks/klipper/klippy/extras
/home/mks/klipper/klippy/extras/gcode_shell_command.py
/home/mks/klipper/klippy/extras/probe.py
/home/mks/klipper/klippy/extras/x_twist_compensation.py
/home/mks/klipper/klippy/gcode.py
/home/mks/klipper/klippy/mcu.py
/home/mks/klipper_config
/home/mks/klipper_config/moonraker.conf
/home/mks/klipper_config/plr.cfg
/home/mks/klipper_config/saved_variables.cfg
/home/mks/moonraker
/home/mks/moonraker/moonraker
/home/mks/moonraker/moonraker/components
/home/mks/moonraker/moonraker/components/file_manager
/home/mks/moonraker/moonraker/components/file_manager/metadata.py
/home/mks/plr.sh
/home/mks/update_moonraker.sh
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/elegoo-fix.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/60-usbmount.rules
/lib/udev/rules.d/99-makerbase-automount.rules
