# 20240430/1.img/var/lib/dpkg/info/network-manager.list
/etc
/etc/NetworkManager
/etc/NetworkManager/NetworkManager.conf
/etc/NetworkManager/dispatcher.d
/etc/NetworkManager/dispatcher.d/01-ifupdown
/etc/init.d
/etc/init.d/network-manager
