# 20240430/1.img/var/lib/dpkg/info/strongswan-starter.list
/etc
/etc/apparmor.d
/etc/apparmor.d/usr.lib.ipsec.stroke
/etc/init.d
/etc/init.d/ipsec
/etc/ipsec.conf
/etc/ipsec.secrets
/etc/strongswan.d
/etc/strongswan.d/charon
/etc/strongswan.d/charon/stroke.conf
/etc/strongswan.d/starter.conf
/var
/var/lib
/var/lib/strongswan
